# .config/systemd/user/gdmflex.service
[Unit]
PartOf=lock.target

[Service]
Type=simple
ExecStart=gdmflexiserver

[Install]
WantedBy=lock.target
WantedBy=sleep.target